# packages/net-im/ejabberd/files/systemd/ejabberd.service
[Unit]
Description=eJabberd
After=network.target

[Service]
Type=forking
User=ejabberd
ExecStart=/usr/sbin/ejabberdctl start
ExecStop=/usr/sbin/ejabberdctl stop

[Install]
WantedBy=multi-user.target
